# homebox/rootfs/etc/services.d/homebox/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Homebox
# This runs Homebox using HA Ingress
# ==============================================================================

# Create data directory if it doesn't exist
data_dir="/data"
mkdir -p "${data_dir}" || bashio::exit.nok "Could not create data directory"

# Set timezone if configured
if bashio::config.has_value 'timezone'; then
    TZ=$(bashio::config 'timezone')
    export TZ
    bashio::log.info "Timezone set to: ${TZ}"
fi

# Verify homebox binary exists and is executable
if [ ! -f /usr/bin/homebox ]; then
    bashio::log.error "Homebox binary not found at /usr/bin/homebox"
    bashio::exit.nok "Homebox binary missing"
fi

if [ ! -x /usr/bin/homebox ]; then
    bashio::log.error "Homebox binary is not executable"
    bashio::exit.nok "Homebox binary not executable"
fi

# Start Homebox
bashio::log.info "Starting Homebox..."
bashio::log.info "Data directory: ${data_dir}"
bashio::log.info "Port: 7745"

# Change to data directory so homebox creates database and files there
cd "${data_dir}"

# Run homebox - use exec to replace shell process with homebox
# This ensures homebox runs as the service process
#
# Key flags:
# --web-host "0.0.0.0" - Bind to all interfaces (required for ingress)
# --web-port 7745 - Port to listen on
# --mode production - Use production mode
# --log-level info - Set logging level
# --log-format text - Use text format for logs (better for addon logs)
# --storage-prefix-path - Where to store uploaded files
# --database-sqlite-path - Full path to SQLite database
# --options-allow-analytics false - Disable analytics
# --options-github-release-check false - Disable GitHub release checks in addon context
exec /usr/bin/homebox \
    --web-host "0.0.0.0" \
    --web-port 7745 \
    --mode production \
    --log-level info \
    --log-format text \
    --storage-prefix-path "${data_dir}/.data" \
    --database-sqlite-path "${data_dir}/.data/homebox.db?_pragma=busy_timeout=999&_pragma=journal_mode=WAL&_fk=1&_time_format=sqlite" \
    --options-allow-analytics false \
    --options-github-release-check false

